# mold_project_gitlab/image/etc/cron.monthly/update_certificates
#!/bin/bash
/opt/mold/scripts/certificate_update
