# mold_project_gitlab/image/etc/cron.hourly/mold_logging
#!/bin/bash
/opt/mold/scripts/logging/hourly_jobs